# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/usr/bin/gunicorn web.server:app -b 0.0.0.0:8080
directory=/var/www/app